# services/receiver.service
# services/receiver.service
[Unit]
Description=GNSS Data Receiver Service

[Service]
ExecStart=/usr/bin/python3 /path/to/project/receiver_service.py <receiver_name>
Restart=always

[Install]
WantedBy=multi-user.target
